# install/src/scot/scot.service
[Unit]
Description=SCOT (Sandia Cyber Omni Tracker)

[Service]
ExecStart=/usr/local/bin/hypnotoad /opt/scot/script/Scot
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
