# deploy/scot3_upstart.conf
description "SCOT3 container"
author "Joshua Maine"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a scot3
end script
